# scripts/com.territory-plan.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.territory-plan.snapshot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npx</string>
    <string>tsx</string>
    <string>scripts/snapshot-data.ts</string>
    <string>--push</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/sierraholstad/Desktop/Interactive Territory Plan Builder/territory-plan</string>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/Users/sierraholstad/Desktop/Interactive Territory Plan Builder/territory-plan/data/snapshots/cron-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sierraholstad/Desktop/Interactive Territory Plan Builder/territory-plan/data/snapshots/cron-stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
